# mysql.conf
[program:mysql]
command=mysqld

autorestart=true
exitcodes=0

stopsignal=TERM
stopasgroup=true
stopwaitsecs=10

stdout_events_enabled=true
stderr_events_enabled=true
